# app/streamvault.service
[Unit]
Description=StreamVault video streaming server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=streamvault
Group=streamvault
WorkingDirectory=/opt/streamvault
ExecStart=/usr/bin/node /opt/streamvault/dist/index.js
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true
ReadWritePaths=/opt/streamvault

[Install]
WantedBy=multi-user.target
